# autorig-online/deploy/autorig.service
[Unit]
Description=AutoRig Online Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/root/autorig-online/backend
Environment="PATH=/root/autorig-online/venv/bin"
# Optional: RESEND_API_KEY, EMAIL_FROM — see deploy/environment.example
EnvironmentFile=-/etc/autorig-online/environment
EnvironmentFile=-/etc/autorig-telegram.env
# Single uvicorn process: one in-process background loop (dispatch/stale). Remote converters are separate HTTP endpoints, not OS workers.
# Use `python -m uvicorn` from the canonical project tree.
ExecStart=/root/autorig-online/venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
# Uvicorn may wait indefinitely on "background tasks" during SIGTERM; cap stop so deploy/restart cannot hang the site.
TimeoutStopSec=45

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/autorig /root/autorig-online/backend/db /root/autorig-online/static/tasks /root/autorig-online/static/glb_cache
PrivateTmp=true

[Install]
WantedBy=multi-user.target

