# autorig-online/deploy/autorig.service.d/no-access-log.conf
# Drop-in: keep service stdout useful by suppressing per-request uvicorn access logs.
# This prevents high-frequency polling endpoints from growing journald into GBs.
[Service]
ExecStart=
ExecStart=/root/autorig-online/venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1 --no-access-log
