# autorig-online/deploy/autorig.service.d/timeout.conf
# Drop-in: do not replace the main autorig.service (ExecStart/User may differ on each host).
# Uvicorn can hang forever on SIGTERM while waiting for "background tasks to complete".
[Service]
TimeoutStopSec=45
